# packages/net/wing/files/etc/uci-defaults/wing
uci delete network.mesh
uci set network.mesh=interface
uci set network.mesh.proto=wing

uci add firewall zone
uci set firewall.@zone[-1].name="mesh"
uci set firewall.@zone[-1].input="ACCEPT"
uci set firewall.@zone[-1].output="ACCEPT"
uci set firewall.@zone[-1].forward="REJECT"
uci set firewall.@zone[-1].masq="1"

uci add firewall forwarding
uci set firewall.@forwarding[-1].src="lan"
uci set firewall.@forwarding[-1].dest="mesh"

uci add firewall forwarding
uci set firewall.@forwarding[-1].src="mesh"
uci set firewall.@forwarding[-1].dest="wan"

