# branches/packages_10.03.2/net/portsentry/files/portsentry.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org
START=99

start() {
	portsentry -atcp
	portsentry -audp
}

stop() {
	killall -9 portsentry
}

