# packages_10.03.1/utils/cmdpad/files/cmdpad.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org
START=93

start() {
    /usr/sbin/cmdpad --quiet > /dev/null &
}
							
stop() {
    killall cmdpad
}
											