# packages/net/parprouted/files/parprouted.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=50

BIN=parprouted
DEFAULT=/etc/default/$BIN

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	$BIN $OPTIONS
}

stop() {
	killall $BIN
}
