# packages_10.03.1/utils/lcd4linux/files/lcd4linux.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org
START=98

NAME=lcd4linux
BIN_F=/usr/bin/$NAME
PID_F=/var/run/lcd4linux.pid
SSD=start-stop-daemon
OPTIONS="-q -p $PID_F -f /etc/lcd4linux.conf"

start() {
	$SSD -S -q -p $PID_F -b -x $BIN_F -- $OPTIONS
}

stop() {
	$SSD -K -q -p $PID_F
}

