# packages/net/rarpd/files/rarpd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=50

SERVICE_DAEMONIZE=1

PROG="rarpd"
COMMAND="/usr/sbin/$PROG"
DEFAULT="/etc/default/$PROG"

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	service_start $COMMAND -n $OPTIONS
}

stop() {
	service_stop $COMMAND
}
