# packages/admin/gkrellmd/files/gkrellmd.init
#!/bin/sh /etc/rc.common

start() {
	/usr/bin/gkrellmd
}

stop() {
	killall gkrellmd
}
