# branches/packages_10.03.1/net/srelay/files/srelay.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=50

start() {
	if [ -f /etc/srelay.conf ];then
		srelay -c /etc/srelay.conf -r -s 
	fi
}

stop() {
	killall srelay
}
