# packages/net/bind/files/named.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

start() {
	if [ -f /etc/bind/named.conf ]; then
		named -c /etc/bind/named.conf
	fi
}

stop() {
	killall named
}
