# branches/packages_10.03.1/net/restund/files/restund.init
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2010 Alfred E. Heggestad
#

START=60
BIN=restund
DEFAULT=/etc/default/$BIN
OPTIONS=

start() {
       [ -f $DEFAULT ] && . $DEFAULT
       $BIN $OPTIONS
}

stop() {
       killall $BIN
}
