# branches/packages_8.09/net/oidentd/files/oidentd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=85
start() {
	/usr/sbin/oidentd -m -f 113
}

stop() {
	killall oidentd
}
