# branches/packages_10.03.2/sound/mpdas/files/mpdas.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009 OpenWrt.org
START=94

start() {
	HOME=/root /usr/bin/mpdas -c /etc/mpdasrc &
}

stop() {
	killall mpdas
}
