# branches/packages_10.03.2/net/seeks/files/seeks.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org
START=99

DAEMON=/usr/bin/seeks

start() {
	echo 'seeks is starting...'
	start-stop-daemon -S -q -x "$DAEMON" -- --daemon
}

stop() {
	echo 'seeks is stopping...'
	start-stop-daemon -K -q -x "$DAEMON"
}

restart() {
	stop
	sleep 2
	start
}
