# branches/packages_10.03.2/net/opennhrp/files/opennhrp.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009 OpenWrt.org
START=50
# Copyright (C) 2009 Jakob Pfeiffer

start() {
	/usr/sbin/opennhrp -d
}

stop() {
	killall opennhrp
}

