# branches/packages_10.03.1/sound/shell-fm/files/shell-fm.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=99
start() {
	[ -d /tmp/.shell-fm ] || mkdir /tmp/.shell-fm
	[ -L /tmp/.shell-fm/shell-fm.rc ] || ln -s /etc/shell-fm.rc /tmp/.shell-fm/shell-fm.rc
	HOME=/tmp shell-fm -d
}

stop() {
	echo "quit" | telnet localhost 54311
}
