# branches/packages_10.03.1/net/ulogd/files/ulogd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=50

BIN=ulogd
DEFAULT=/etc/default/$BIN
LOG_D=/var/log

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	$BIN $OPTIONS
}

stop() {
	killall $BIN
}
