# packages/net/privoxy/files/privoxy.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
conf=/etc/privoxy/config

start() {
        /usr/sbin/privoxy $conf
}

stop() {
        killall privoxy
}
