# branches/packages_10.03.1/net/ndyndns/files/ndyndns.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org

START=80
DEFAULT=/etc/default/ndyndns

start() {
	[ -f $DEFAULT ] && . $DEFAULT

	. /etc/functions.sh
	include /lib/network
	scan_interfaces
	/usr/sbin/ndyndns -c $CHROOT -f $CONFIG -i $(config_get wan ifname) $OPTIONS
}

stop() {
	killall ndyndns
}
