# packages/net/openntpd/files/ntpd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=60

boot() {
	return 0
}

start() {
	mkdir -p `awk -F: '/^ntp:/{print $6}' /etc/passwd`
	mkdir -p /var/db/
	/usr/sbin/ntpd -s
}

stop() {
	killall ntpd
}
