# packages/utils/setserial/files/serial.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

start() {
	/usr/sbin/setserial /dev/tts/1 irq 3
}
