# branches/packages_10.03.1/net/unfs3/files/unfs3.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=50

start() {
	/usr/sbin/unfsd
}

stop() {
	killall unfsd
}

restart() {
	stop
	sleep 1
	start
}
