# packages/net/nfs-server/files/nfsd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=50

start() {
	/usr/sbin/portmap
	/usr/sbin/rpc.mountd -r
	/usr/sbin/rpc.nfsd
}

stop() {
	killall rpc.nfsd
	killall rpc.mountd
	killall portmap
}
