# packages/net/pure-ftpd/files/pure-ftpd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

BIN=pure-ftpd
DEFAULT=/etc/default/$BIN
LOG_D=/var/log/$BIN
RUN_D=/var/run
PID_F=$RUN_D/$BIN.pid
OPTIONS="-A -B -c 10 -C 4 -E -lunix -U 133:022 -g $PID_F -y 3:4"
[ -f $DEFAULT ] && . $DEFAULT

start() {
        mkdir -p $LOG_D
        mkdir -p $RUN_D
        $BIN $OPTIONS
}

stop() {
        [ -f $PID_F ] && kill $(cat $PID_F)
}
