# branches/packages_10.03.2/net/vpzone/files/vpzone.init
#!/bin/sh /etc/rc.common

START=50
DAEMON=vpzone
PID_F=/var/run/$DAEMON.pid

start() {
	start-stop-daemon -p $PID_F -x $DAEMON -S -b -m
}

stop() {
	start-stop-daemon -p $PID_F -x $DAEMON -K
}
