# packages/Xorg/wm/fluxbox/files/uci-defaults/x11
uci set x11.@desktop[0].autostart=enlightenment_start