# packages/phone/paroli/files/paroli.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org
START=99

start() {
    (sleep 10 && DISPLAY=":0" LD_PRELOAD=/usr/lib/libX11.so /usr/bin/paroli-launcher) & # sleep to wait until the xserver is started | LD_PRELOAD to workaround some runtime linking problems - needs to be fixed
}
    
stop() {
    killall paroli-launcher
}
