# packages/net/portmap/files/portmap.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=40
BIN=portmap
DEFAULT=/etc/default/$BIN

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	$BIN $OPTIONS
}

stop() {
	killall $BIN
}

