# packages/net/portmap/files/portmap.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

BIN=portmap
DEFAULT=/etc/default/$BIN
[ -f $DEFAULT ] && . $DEFAULT

start() {
	$BIN $OPTIONS
}

stop() {
	killall $BIN
}

