# packages/net/frickin/files/frickin.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=50

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

DEFAULT=/etc/default/frickin

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	service_start /usr/sbin/frickin $OPTIONS
}

stop() {
	service_stop /usr/sbin/frickin
}
