# branches/packages_10.03.2/net/aprx/files/aprx.init
#!/bin/sh /etc/rc.common
#Init script for aprx

START=99
APRX_BIN="/usr/sbin/aprx"
APRX_LOG_DIR="/var/log/aprx"

system_config() {
	local cfg="$1"

	config_get hostname "$cfg" hostname
}


start() {
	[ -x "$APRX_BIN" ] || return 1
	[ -d "$APRX_LOG_DIR" ] || mkdir -p $APRX_LOG_DIR
	$APRX_BIN
	}

stop() {
	killall aprx
}
