# packages/net/ptunnel/files/ptunnel.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org
START=70
BIN=ptunnel

start() {
    config_load ptunnel
    config_foreach start_ptunnel ptunnel
}

start_ptunnel() {
    local section="$1"
    
    include /lib/network
    scan_interfaces
    
    iptables -I OUTPUT 1 -p icmp -j ACCEPT
    config_get interface "$section" interface

    $BIN -c "$interface" > /dev/null&
}

stop() {
        ps aux | grep ptunnel |  awk '{print }' | \
            xargs kill 2> /dev/null
}
									