# branches/packages_8.09/net/stun/files/stund.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009 OpenWrt.org
START=70

start() {
	config_load "6stund"

	config_get PRIMIP primaryip
	config_get SECIP secondaryip
	config_get PRIMPORT primaryport
	config_get SECPORT secondaryport

	/usr/sbin/stund -b -h $PRIMIP -a $SECIP -p $PRIMPORT -o $SECPORT
}

stop() {
	killall -9 stund
}
