# packages/net/rp-pppoe/files/pppoe-client.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

start() {
	pppoe-start
}

stop() {
	pppoe-stop
}

