# packages/net/rp-pppoe/files/pppoe-server.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

DEFAULT=/etc/default/pppoe-server
[ -f $DEFAULT ] && . $DEFAULT

start() {
	pppoe-server $OPTIONS
}

stop() {
	killall pppoe-server
}
