# packages/net/rp-pppoe/files/pppoe-relay.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

DEFAULT=/etc/default/pppoe-relay
[ -f $DEFAULT ] && . $DEFAULT

start() {
	pppoe-relay $OPTIONS
}

stop() {
	killall pppoe-relay
}

