# packages/ipv6/miredo/files/miredo-server.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=50

PID_F=/var/run/miredo-server.pid

start() {
	insmod ipv6 2>/dev/null
	insmod tun 2>/dev/null
	miredo-server -c /etc/miredo-server.conf
}

stop() {
	[ -f $PID_F ] && kill $(cat $PID_F)
}
