# packages/ipv6/miredo/files/miredo.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=50

PID_F=/var/run/miredo.pid

start() {
	insmod ipv6 2>/dev/null
	insmod tun 2>/dev/null
	miredo -c /etc/miredo.conf
}

stop() {
	[ -f $PID_F ] && kill $(cat $PID_F)
}
