# branches/packages_10.03.1/net/rp-l2tp/files/rp-l2tpd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=60
DEFAULT=/etc/default/rp-l2tpd
RUN_D=/var/run
[ -f $DEFAULT ] && . $DEFAULT

start() {
	[ -d $RUN_D ] || mkdir -p $RUN_D
	l2tpd $OPTIONS
}

stop() {
	killall l2tpd
}

