# branches/packages_10.03.1/net/rp-pppoe/files/pppoe-server.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=50

DEFAULT=/etc/default/pppoe-server

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	pppoe-server $OPTIONS
}

stop() {
	killall pppoe-server
}
