# packages/net/batman-advanced/files/etc/init.d/batman-adv-kernelland
#!/bin/sh /etc/rc.common
START=90

start () {
	interfaces=$(uci get batman-adv-kernelland.general.interface)
	if [ "$interfaces" = "" ]; then
		echo $1 Error, you must specify at least a network interface
		exit
	fi
	
	echo $interfaces > /proc/net/batman-adv/interfaces 
	
	originator_interval=$(uci get batman-adv-kernelland.general.originator_interval)
	log_level=$(uci get batman-adv-kernelland.general.log_level)

	if [ $originator_interval ]; then
		echo $originator_interval > /proc/net/batman-adv/orig_interval 
	fi
	
	if [ $log_level ]; then
		echo $log_level > /sys/module/batman_adv/parameters/debug 2>&-
	fi
}

stop () {
	echo "" > /proc/net/batman-adv/interfaces
}

        
