# packages/net/netatalk/files/afpd.init
#!/bin/sh /etc/rc.common

START=70

MAXCONS="7"
CONFIGFILE="/etc/netatalk/afpd.conf"
PIDFILE="/var/run/afpd.pid"

start()
{
	/usr/bin/cnid_metad
	/usr/bin/afpd -F ${CONFIGFILE} -P ${PIDFILE} -c ${MAXCONS}
}

stop()
{
	killall afpd
	killall cnid_metad
}

reload()
{
	kill -s HUP `cat ${PIDFILE}`
	start
}
