# packages/utils/dbus/files/dbus.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org

SSD=start-stop-daemon
PIDF=/var/run/dbus.pid
PROG=/usr/sbin/dbus-daemon
ARGS=

START=60
start() {
	mkdir -p /var/run/dbus
	$SSD -S -p $PIDF -q -x $PROG -- $ARGS 
}

stop() {
	$SSD -K -p $PIDF -q
}

status() {
        ps | grep -q "^ *$(cat $PIDF 2>/dev/null) *.*$PROG"
}

