# packages/lang/php5/files/php.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2010 OpenWrt.org

START=50

export PHP_FCGI_CHILDREN=''
PORT=1026
BIN=/usr/bin/php-fcgi

start() {
	$BIN -b $PORT &
}

stop() {
	kill `pidof php-fcgi`
}
