# packages/net/nfs-kernel-server/files/nfsd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=60
STOP=60

NFS_D=/var/lib/nfs
LOCK_D=/var/lib/nfs/sm

start() {
	mkdir -p $NFS_D
	mkdir -p $LOCK_D
	touch $NFS_D/rmtab

	/usr/sbin/rpc.statd
	/usr/sbin/exportfs -r
	/usr/sbin/rpc.nfsd
	/usr/sbin/rpc.mountd
}

stop() {
	killall rpc.mountd 2> /dev/null
	killall -9 rpc.nfsd 2> /dev/null
	/usr/sbin/exportfs -au
	killall rpc.statd 2> /dev/null
}

reload() {
	/usr/sbin/exportfs -r
}

