# branches/packages_10.03.1/net/netatalk/files/afpd.init
#!/bin/sh /etc/rc.common

START=70

MAXCONS="7"
CONFIGFILE="/etc/netatalk/afpd.conf"
PIDFILE="/var/run/afpd.pid"

start()
{
	/usr/sbin/cnid_metad
	/usr/sbin/afpd -F ${CONFIGFILE} -P ${PIDFILE} -c ${MAXCONS}
}

stop()
{
	kill `cat "${PIDFILE}"`
	killall cnid_metad
}

reload()
{
	kill -s HUP `cat "${PIDFILE}"`
}
